# overlay/system/etc/profile.d/pinephone-camera-orientation.sh
# Set an environment variable so that Qt knows the camera orientation.
export QT_GSTREAMER_UBPORTS_CAM_ORIENTATION=270
